# .conf/supervisor-app.conf
[program:uwsgi]
command = uwsgi --ini /etc/uwsgi/sites/app.ini

[program:nginx]
<<<<<<< HEAD
command = nginx
=======
command= nginx

>>>>>>> bc95b67aef65e6f984d661d45e80409205842039
